# scripts/com.superdownloads.health-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.superdownloads.health-check</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/supermac/Desktop/DEV/01_Projects/SUPER-DOWNLOADS/scripts/platform-health-notify.sh</string>
  </array>

  <!-- Daily at 10:00. launchd runs a missed job at next wake if the Mac was asleep. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/supermac/Library/Logs/super-downloads-health.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/supermac/Library/Logs/super-downloads-health.err.log</string>
</dict>
</plist>
